# conf/supervisor/laravel-queue.ini
[program:laravel-queue]
user=www-data
command=php /var/www/html/artisan queue:work --sleep=3 --tries=5 --delay=3
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/dev/stdout
